# amn-via-socks.service
[Unit]
Description=Redirect container traffic through tun2socks
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/root/amn-via-socks.sh
ExecStop=/root/amn-via-socks.sh stop
Restart=always
RestartSec=5s
User=root
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
